# blueprints/cloud-operations/uptime-monitoring-tool/data/config-sync.conf
[program:configsync]
directory=/opt/config_sync
command=/opt/config_sync/env/bin/python /opt/config_sync/source/main.py
autostart=true
autorestart=true
user=net-mon-config-sync
# Environment variables ensure that the application runs inside of the configured virtualenv.
environment=VIRTUAL_ENV="/opt/config_sync/env",PATH="/opt/config_sync/env/bin",HOME="/home/net-mon-config-sync",USER="net-mon-config-sync"
stdout_logfile=syslog
stderr_logfile=syslog
